# infracode/ansible/templates/etc_supervisor_conf.d_load-balancer.conf
[program:load-balancer]
command={{lb_daemon_invocation}}
user=root
group=root
autostart=true
autorestart=true
stderr_logfile=syslog
stdout_logfile=syslog
